# templates/launchd/brain-auto-extract.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.{{USERNAME}}.brain-auto-extract</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>{{BRAIN_DIR}}/bin/auto-extract.sh</string>
    </array>

    <!-- Safety-net poll every 5 minutes. Day-to-day we run via WatchPaths. -->
    <key>StartInterval</key>
    <integer>300</integer>

    <!-- Event-driven. The vault root + Claude project dir + Cursor agent
         transcript dir are watched so any Obsidian save or session message
         triggers a run. Cursor was added 2026-04-20 — without it, Cursor
         session captures had to wait for the 5-min safety poll. -->
    <key>WatchPaths</key>
    <array>
        <string>{{BRAIN_DIR}}</string>
        <string>{{BRAIN_DIR}}/raw</string>
        <string>{{HOME}}/.claude/projects</string>
        <string>{{HOME}}/.cursor/projects</string>
    </array>

    <!-- 1 s throttle: pure-Python ingest is cheap (~50 ms idle). -->
    <key>ThrottleInterval</key>
    <integer>1</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>{{BRAIN_DIR}}/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>{{BRAIN_DIR}}/logs/launchd-stderr.log</string>

    <key>WorkingDirectory</key>
    <string>{{BRAIN_DIR}}</string>

    <!-- Inject BRAIN_DIR so brain.config picks up the user's chosen vault.
         launchd doesn't read the user's shell rc — without this, the
         Python subprocess would fall back to ~/.brain. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>BRAIN_DIR</key>
        <string>{{BRAIN_DIR}}</string>
    </dict>
</dict>
</plist>
