# templates/launchd/brain-auto-extract.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.{{USERNAME}}.brain-auto-extract</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>{{BRAIN_DIR}}/bin/auto-extract.sh</string>
    </array>

    <!-- Safety-net poll every 5 minutes. Day-to-day we run via WatchPaths. -->
    <key>StartInterval</key>
    <integer>300</integer>

    <!-- Event-driven. Watch INPUT sources only — never the vault root.
         auto-extract writes to logs/, entities/_MOC.md, .harvested,
         .brain.db, .extract.lock.d/, all inside {{BRAIN_DIR}}. Watching
         the root caused a self-trigger feedback loop (~1 run/sec, 125k
         runs in <24h, log to 57 MB) that hammered disk I/O hard enough
         to block pyenv-rehash and freeze new terminal startup for 60s
         (incident 2026-04-28). The 5-min StartInterval above is the
         safety net for anything that doesn't fall under an input-only
         watch. -->
    <key>WatchPaths</key>
    <array>
        <string>{{BRAIN_DIR}}/raw</string>
        <string>{{HOME}}/.claude/projects</string>
        <string>{{HOME}}/.cursor/projects</string>
    </array>

    <!-- 30 s throttle. Bursty session writes (every keystroke updates the
         transcript jsonl) used to fire back-to-back at 1 s, generating
         needless Python interpreter spawns. 30 s is short enough that
         the user perceives ingest as "live" but long enough that
         coalescing cuts CPU pressure ~30×. -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>{{BRAIN_DIR}}/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>{{BRAIN_DIR}}/logs/launchd-stderr.log</string>

    <key>WorkingDirectory</key>
    <string>{{BRAIN_DIR}}</string>

    <!-- Inject BRAIN_DIR so brain.config picks up the user's chosen vault.
         launchd doesn't read the user's shell rc — without this, the
         Python subprocess would fall back to ~/.brain. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>BRAIN_DIR</key>
        <string>{{BRAIN_DIR}}</string>
    </dict>
</dict>
</plist>
